# packaging/ccsockd/ccsockd.service
[Unit]
Description=Headless Claude Code daemon
After=default.target

[Service]
ExecStart=%h/.local/bin/ccsockd
Restart=on-failure
RestartSec=2s

[Install]
WantedBy=default.target
