# packaging/blemees-agentd/com.blemees.blemees-agentd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.blemees.blemees-agentd</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/env</string>
        <string>blemees-agentd</string>
    </array>

    <!--
        If `claude` or `codex` aren't on launchd's PATH, point at them
        explicitly:
            <key>EnvironmentVariables</key>
            <dict>
                <key>BLEMEES_AGENTD_CLAUDE</key>
                <string>/Users/you/.local/bin/claude</string>
                <key>BLEMEES_AGENTD_CODEX</key>
                <string>/Users/you/.local/bin/codex</string>
            </dict>
        Both default to the bare binary name resolved through PATH.
    -->

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>5</integer>

    <key>StandardErrorPath</key>
    <string>/tmp/blemees-agentd.err.log</string>

    <key>StandardOutPath</key>
    <string>/tmp/blemees-agentd.out.log</string>
</dict>
</plist>
